# dotfiles/opentabletdriver/systemd/opentabletdriver.service
Unit]
Description=OpenTabletDriver Service
Documentation=https://github.com/OpenTabletDriver/OpenTabletDriver
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/OpenTabletDriver
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
Environment=XDG_CONFIG_HOME=%h/.config
Environment=XDG_DATA_HOME=%h/.local/share
Environment=XDG_CACHE_HOME=%h/.cache
# Wayland support
Environment=XDG_SESSION_TYPE=wayland
Environment=QT_QPA_PLATFORM=wayland
Environment=GDK_BACKEND=wayland
# Ensure proper permissions
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/OpenTabletDriver
ReadWritePaths=%h/.local/share/OpenTabletDriver
ReadWritePaths=/dev/input
ReadWritePaths=/dev/usb

Install]
WantedBy=graphical-session.target 